# systemd_units/selfstorage_reminders.service
[Unit]
Description=Send rent reminders

[Service]
WorkingDirectory=/opt/SelfStorage/
ExecStart=venv/bin/python manage.py send_rent_reminders
