# systemd_units/selfstorage_clearsessions.timer
[Unit]
Description=Timer for clearsessions of SelfStorage site

[Timer]
OnBootSec=60
OnUnitActiveSec=1d

[Install]
WantedBy=multi-user.target
