# systemd_units/selfstorage_notices.service
[Unit]
Description=Send expired rent notices

[Service]
WorkingDirectory=/opt/SelfStorage/
ExecStart=venv/bin/python manage.py send_expired_notices
