# systemd_units/selfstorage_reminders.timer
[Unit]
Description=Timer for rent reminders
PartOf=selfstorage.target

[Timer]
OnCalendar=*-*-* 12:00:00

[Install]
WantedBy=multi-user.target

