# systemd_units/selfstorage_reminders.service
[Unit]
Description=Send rent reminders

[Service]
ExecStart=/opt/SelfStorage/venv/bin/python manage.py send_rent_reminders
